# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx 
user=root
autostart=true

[program:tomcat]
command=/bin/bash /tomcat_run.sh
