# supervisord.conf
[supervisord]
nodaemon=true

[program:my_init]
command=/sbin/my_init

[program:nginx]
command=service nginx start

[program:tomcat]
command=service tomcat7 start
